# apache-exporter.conf
start on startup

script
  /root/go/bin/apache_exporter
end script
